# cape/etc/systemd/system/cape-suricata.service
[Unit]
Description=CAPE Suricata IDS/IDP daemon
After=network.target
Requires=network.target

[Service]
Type=forking
ExecStartPre=/bin/rm -f /var/run/cape-suricata.pid
ExecStart=/usr/bin/suricata -D -c /etc/suricata/cape-suricata.yaml -k none --unix-socket --user cape --group cape
ExecStartPost=/bin/chown cape:cape -R /var/run/cape-suricata
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill $MAINPID
PrivateTmp=no
InaccessibleDirectories=/home /root
ReadOnlyDirectories=/boot /usr /etc
User=root

[Install]
WantedBy=multi-user.target

