# ghidra/etc/systemd/system/ghidra.service
[Unit]
Description=Ghidra server

[Service]
Type=simple
User=ghidra
Group=ghidra
ExecStart=/opt/ghidra/server/ghidraSvr console
Restart=always
RestartSec=5m

[Install]
WantedBy=multi-user.target
