# ghidra/etc/profile.d/ghidra.sh
export PATH=\$PATH:/usr/lib/jvm/jre-11-openjdk/
export JAVA_HOME=/usr/lib/jvm/jre-11-openjdk/
export GHIDRA_HOME=/opt/ghidra/
